# ptit-firewall/ptit-firewall/ptit-crontab/victim/_system/etc/cron.d/cronjob_victim
@reboot victim /usr/bin/cronjob_victim.sh &> /dev/null
* * * * * victim /usr/bin/cronjob_victim.sh &> /dev/null
